# systemd/perpprompt-dreaming.service
[Unit]
Description=perpPROMPT DREAMING cycle — the bird sleeps and surfaces fragments
Documentation=https://github.com/ASIXicle/perpPROMPT
After=network-online.target llama-server.service llama-server-embedding.service
Wants=network-online.target
Requires=llama-server.service llama-server-embedding.service

[Service]
Type=oneshot
User=perpbot
Group=perpbot
WorkingDirectory=/opt/perpbot

EnvironmentFile=-/opt/perpbot/config/bird_env

# DREAMING cycles weight-alternate between the utility variant (dream.md —
# identity + dreaming, dreamer-autonomy storage) and the free variant
# (dream.free.md — identity-free pure generation, wrapper auto-stores).
# The split is controlled by DREAM_FREE_WEIGHT in config/bird_env: the
# probability [0.0-1.0] that a given cycle runs free (default 0.0 =
# utility-only). dream.py rolls the die per invocation — no --free flag
# needed here, no sibling unit. Tune the dial by editing bird_env, not
# this file.
#
# Phase E note: DREAMING stays timer-disabled for the first 5 days
# (THINKING-only window per dry-run-evaluation.md gates); Holden enables
# the timer after that. Set DREAM_FREE_WEIGHT > 0 in bird_env when ready
# to mix free dreams into the rotation.
ExecStart=/opt/perpbot/venv/bin/python -m src.dream --bird-name ${BIRD_NAME}

# 15-minute hard cap. DREAMING worst case is 5 turns × ~2min = 10min on
# Sandy Bridge. 15min catches runaways without killing legitimate cycles.
TimeoutStartSec=900

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ReadWritePaths=/opt/perpbot
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
