# systemd/perpprompt-dreaming.timer
[Unit]
Description=perpPROMPT DREAMING cadence — bird sleeps every 4-8h (~6h avg)
Documentation=https://github.com/ASIXicle/perpPROMPT

[Timer]
# First fire 4 hours after boot — significantly later than THINKING's 15min
# offset, so the first hours after boot are THINKING-only. This intentionally
# offsets DREAMING from THINKING's cadence so they don't pile up.
OnBootSec=4h

# 4h base interval + up to 4h randomized delay = 4-8h between cycles, ~6h
# average (~3-4 dreams/day). systemd jitter only ADDS time, so a 4h base + 4h
# jitter is how we render "6h ± 2h": some cycles fire fast (~4h, catching
# momentum), some slow (~8h, letting fragments decay). Bumped from 10-12h to
# feed more DREAM-FREE content to the Bluesky feed (2026-06-02 chorus round).
OnUnitActiveSec=4h
RandomizedDelaySec=4h

Persistent=true

Unit=perpprompt-dreaming.service

[Install]
WantedBy=timers.target
