# systemd/perpprompt-dreaming.service
[Unit]
Description=perpPROMPT DREAMING cycle — the bird sleeps and surfaces fragments
Documentation=https://github.com/ASIXicle/perpPROMPT
After=network-online.target llama-server.service llama-server-embedding.service
Wants=network-online.target
Requires=llama-server.service llama-server-embedding.service

[Service]
Type=oneshot
User=perpbot
Group=perpbot
WorkingDirectory=/opt/perpbot

EnvironmentFile=-/opt/perpbot/config/bird_env

# DREAMING cycles randomly alternate between the utility variant (dream.md)
# and the artistic variant (dream.free.md). For Phase E first 5 days,
# DREAMING is DISABLED at the timer level — Holden enables only after the
# THINKING-only window completes per dry-run-evaluation.md gates.
# When enabled, alternation can be configured by toggling --free between
# this service file's ExecStart and a sibling. For initial deployment,
# omit --free (utility variant only); add --free toggling in Phase E+.
ExecStart=/opt/perpbot/venv/bin/python -m src.dream --bird-name ${BIRD_NAME}

# 15-minute hard cap. DREAMING worst case is 5 turns × ~2min = 10min on
# Sandy Bridge. 15min catches runaways without killing legitimate cycles.
TimeoutStartSec=900

StandardOutput=journal
StandardError=journal

ProtectSystem=strict
ReadWritePaths=/opt/perpbot
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
